# dist/freebsd/net/reflector/files/reflector.in
#!/bin/sh
#
# PROVIDE: reflector
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following to /etc/rc.conf to enable the reflector:
#
#   reflector_enable="YES"
#
# Optional:
#
#   reflector_config="/path/to/reflector.toml"   # default %%PREFIX%%/etc/reflector.toml
#   reflector_flags=""                           # extra arguments

. /etc/rc.subr

name="reflector"
rcvar="reflector_enable"
desc="Reflect link-local service discovery between interfaces"

load_rc_config $name

: ${reflector_enable:="NO"}
: ${reflector_config:="%%PREFIX%%/etc/reflector.toml"}

pidfile="/var/run/${name}.pid"
procname="%%PREFIX%%/bin/${name}"

# The daemon runs in the foreground and logs to stderr, so daemon(8) backgrounds it and -S routes that
# stderr into syslog under our own tag rather than /dev/null.
command="/usr/sbin/daemon"
command_args="-S -T ${name} -p ${pidfile} ${procname} ${reflector_config}"

start_precmd="${name}_precmd"

reflector_precmd()
{
    if [ ! -f "${reflector_config}" ]; then
        err 1 "${reflector_config} does not exist"
    fi

    # Refuse to start on a configuration the daemon would reject, so the failure is one clear message
    # here rather than a service that appears to start and is gone a moment later. --check-config parses
    # and validates only: it opens no interface and needs no privileges.
    if ! ${procname} --check-config "${reflector_config}" > /dev/null 2>&1; then
        ${procname} --check-config "${reflector_config}" 2>&1 | while read -r line; do
            warn "${line}"
        done
        err 1 "refusing to start: ${reflector_config} is not valid"
    fi
}

run_rc_command "$1"
